# overlay/system/etc/init/dt2w.conf
# Set permissions for DT2W

description "Set permissions for DT2W"

start on android

task

script
   chmod 0660 /proc/nvt_wake_gesture
   chown system system /proc/nvt_wake_gesture
end script
