# overlay/system/etc/init/widi.conf
# Enable WiFi display

description "Enable WiFi Display"

start on android

task

script
   setprop ubuntu.widi.supported 1
end script
